# files/mrtg.cron
*/5 * * * * nginx /usr/bin/mrtg /etc/mrtg/mrtg.cfg || true
