# files/sv-crond.ini
[program:crond]
command=/usr/sbin/crond -f -L /proc/self/fd/1 -l debug
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/crond.err.log
stdout_logfile=/var/log/supervisor/crond.out.log
