# files/sv-lighttpd.ini
[program:lighttpd]
command=/usr/sbin/lighttpd -f /etc/lighttpd/lighttpd.conf -D
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/lighttpd.err.log
stdout_logfile=/var/log/supervisor/lighttpd.out.log
